# build/agent/agent.service
[Unit]
Description=paas-operator
Wants=network-online.target
After=network.target

[Service]
Type=simple
RemainAfterExit=yes
ExecStart=/usr/local/bin/agent
ExecReload=/usr/local/bin/agent

[Install]
WantedBy=multi-user.target
